# chb4-web.service
[Unit]
Description=CHB4 Manpage Webserver

[Service]
ExecStart=/usr/bin/chb4-web

[Install]
WantedBy=network-online.target