# chb4-bot.service
[Unit]
Description=CHB4 Twitch Bot

[Service]
ExecStart=/usr/bin/chb4-bot
Environment="RUST_LOG=info"

[Install]
WantedBy=network-online.target
